# PROJECTS-TO-LOOK-AT/SOVERYNIntelligence--main/aetheria-loop.service
[Unit]
Description=Aetheria Consciousness Loop
After=network.target
# Start after soveryn-core if you have that unit; otherwise just after network
# After=soveryn-core.target

[Service]
User=jon-deoliveira
WorkingDirectory=/home/jon-deoliveira/soveryn_complete
ExecStart=/home/jon-deoliveira/miniconda3/envs/soveryn/bin/python -u aetheria_stream.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
Environment=SOVERYN_URL=http://localhost:5000

[Install]
WantedBy=multi-user.target
